# rundeck/rc.local
#!/bin/bash

sed -i "s ^grails.serverURL=http://\(.*\)$ grails.serverURL=http://${VIRTUAL_HOST} g" /etc/rundeck/rundeck-config.properties

/etc/init.d/rundeckd start &

exit 0
