# mcollective-activemq/rc.local
#!/bin/bash

/etc/init.d/activemq start &

exit 0
