# jekyll/rc.local
#!/bin/bash

cd /var/www/jekyll

jekyll serve

exit 0
