# rethinkdb/rc.local
#!/bin/bash

/etc/init.d/rethinkdb start &

exit 0
