# ka-lite/rc.local
#!/bin/bash

service kalite start

exit 0
