# rabbitmq/rc.local
#!/bin/bash

/usr/local/bin/rabbitmq-start &

exit 0
