# mongodb/rc.local
#!/bin/bash

sudo -u mongodb '/usr/bin/mongod' &

exit 0
