# flask/rc.local
#!/bin/bash

su - flask -c 'cd /var/www && ./run.py' &

exit 0