# base/my_exec.cron
#!/bin/bash

if    [ "${ENABLE_CRON}" == "false" ] \
   || [ -z "${DISABLE_ALL}" ]; then
  echo "cron disabled"
  chmod -x /etc/service/cron/run > /dev/null 2>&1
fi
