# base/my_exec.cron
#!/bin/bash

if    [ "${ENABLE_CRON}" == "false" ] \
   || [ -z "${DISABLE_ALL}" ]; then
  rm -rf /etc/service/cron > /dev/null 2>&1
fi
